# todoist-popup/com.todoist-popup.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.todoist-popup.hourly</string>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
  </array>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>WorkingDirectory</key>
  <string>/Users/user/l/todoist-popup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/user/l/todoist-popup/popup-todoist</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/Users/user/Library/Logs/todoist-popup.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/user/Library/Logs/todoist-popup.error.log</string>
</dict>
</plist>
